# config/supervisor/claude-telegram.conf
[program:claude-telegram]
command=/home/claude/.claude-super-setup/scripts/vps-telegram-supervisor.sh
user=claude
directory=/home/claude
autostart=true
autorestart=true
startsecs=10
startretries=999
stopwaitsecs=30
stdout_logfile=/home/claude/.claude/logs/telegram-supervisor-stdout.log
stderr_logfile=/home/claude/.claude/logs/telegram-supervisor-stderr.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=3
environment=HOME="/home/claude",PATH="/home/claude/.nvm/versions/node/v22.22.2/bin:/home/claude/.bun/bin:/usr/local/bin:/usr/bin:/bin",NVM_DIR="/home/claude/.nvm"
