# config/systemd/claude-telegram@.service
[Unit]
Description=Claude Code Telegram Listener (24/7)
Documentation=https://github.com/Calebmambwe/claude-super-setup
After=network-online.target
Wants=network-online.target
# Reset restart counter every 10 minutes — prevents permanent failure
StartLimitIntervalSec=600
StartLimitBurst=20

[Service]
Type=simple
User=%i
WorkingDirectory=%h/.claude-super-setup
ExecStart=/usr/local/bin/claude --dangerously-skip-permissions --channels plugin:telegram@claude-plugins-official
StandardOutput=journal
StandardError=journal
# ALWAYS restart — whether clean exit, crash, or signal kill
Restart=always
RestartSec=5
# Watchdog: process must notify every 120s or systemd kills and restarts it
WatchdogSec=120
# Environment — ANTHROPIC_API_KEY or CLAUDE_CODE_OAUTH_TOKEN must be set
# Use systemctl edit claude-telegram@$USER to add secrets
EnvironmentFile=-%h/.claude/env
Environment=HOME=%h
# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.claude %h/.claude-super-setup %h/claude_super_setup
CapabilityBoundingSet=
# Resource limits
LimitNOFILE=65536
MemoryMax=2G

[Install]
WantedBy=multi-user.target
